# misc/websocket-notify.service
[Unit]
Description=websocket-notify

[Service]
ExecStart=/usr/local/bin/websocket-notify

[Install]
WantedBy=default.target
